# setup/claudeclaw@.service
[Unit]
Description=ClaudeClaw Agent: %i
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=forking
WorkingDirectory=%h/workspace
ExecStart=/bin/bash -lc '\
  AGENT=%i; \
  if [ "$AGENT" = "primary" ]; then \
    AGENT_DIR="%h/workspace"; \
    EXTRA_FLAGS=""; \
  else \
    AGENT_DIR="%h/workspace/agents/$AGENT"; \
    EXTRA_FLAGS="--add-dir %h/workspace"; \
  fi; \
  tmux -L "agent-$AGENT" new-session -d -s "agent-$AGENT" \
    "cd $AGENT_DIR && claude $EXTRA_FLAGS --channels plugin:telegram@claude-plugins-official"'
ExecStop=/usr/bin/tmux -L agent-%i kill-session -t agent-%i
Restart=on-failure
RestartSec=30

# Environment
EnvironmentFile=-%h/workspace/.env.agents
Environment=HOME=%h
Environment=PATH=%h/.bun/bin:%h/.npm-global/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
